# conf/etc/services.d/01-cups/run
#!/bin/sh
NAME="CUPS" 
COMMAND="cupsd -f" 
WAIT="/tmp/cups"
exec service-handler "$NAME" "$COMMAND" "$WAIT"